# system_files/usr/lib/systemd/user/cosmic-ext-alternative-startup.service
[Unit]
Description=COSMIC Alternative Startup Service
PartOf=graphical-session.target
After=cosmic-session.service
Requisite=cosmic-session.service
ConditionEnvironment=XDG_CURRENT_DESKTOP=niri

[Service]
Environment=XDG_SESSION_TYPE=wayland
ExecStart=/usr/bin/cosmic-ext-alternative-startup
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
