# system_files/usr/lib/systemd/user/xdg-desktop-portal-gtk.service
[Unit]
Description=Portal service (GTK/GNOME implementation)
PartOf=graphical-session.target
After=niri.service
Requires=niri.service

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.gtk
ExecStart=/usr/libexec/xdg-desktop-portal-gtk

[Install]
WantedBy=graphical-session.target
